# saltprophet.service
[Unit]
Description=Pi-Autopilot Reddit-to-Gumroad Pipeline
After=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/pi-autopilot
Environment="PATH=/opt/pi-autopilot/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/pi-autopilot/.env
ExecStart=/opt/pi-autopilot/venv/bin/python /opt/pi-autopilot/main.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-autopilot

[Install]
WantedBy=multi-user.target
